# ims-integration.service
[Unit]
Description=IMS Integration Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
Group=docker
WorkingDirectory=/opt/ims-integration
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Pull latest image and start container
ExecStartPre=/usr/bin/docker-compose -f docker-compose.prod.yml pull
ExecStart=/usr/bin/docker-compose -f docker-compose.prod.yml up

# Stop gracefully
ExecStop=/usr/bin/docker-compose -f docker-compose.prod.yml down

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target